# systemd/jwb-rpi@.service
[Unit]
Description=JW Broadcasting player

[Service]
ExecStart=/usr/local/bin/jwb-rpi
User=%i
Restart=always

[Install]
WantedBy=multi-user.target
