# systemd/jwb-index@.timer
[Timer]
OnCalendar=daily
Persistent=yes

[Install]
WantedBy=multi-user.target
